# deploy/systemd/dhan-live-trader.timer
# D2: systemd timer unit for dhan-live-trader
#
# Triggers the service at 08:15 AM IST daily.
# Persistent=true fires on next boot if EC2 was off at 08:15.
# AccuracySec=1us disables the default 60s randomized delay.
#
# Install (with service unit):
#   sudo cp deploy/systemd/dhan-live-trader.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable dhan-live-trader.timer
#
# Check next trigger:
#   systemctl list-timers dhan-live-trader.timer

[Unit]
Description=Daily Dhan Live Trader Startup

[Timer]
OnCalendar=*-*-* 08:15:00 Asia/Kolkata
AccuracySec=1us
Persistent=true
WakeSystem=false

[Install]
WantedBy=timers.target
