# deploy/systemd/tickvault.service
# =============================================================================
# tickvault — systemd Service Unit
# =============================================================================
# Auto-restarts the trading app on crash. Watchdog kills hung processes.
#
# Install:
#   sudo cp deploy/systemd/tickvault.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tickvault
#   sudo systemctl start tickvault
#
# Check:
#   systemctl status tickvault
#   journalctl -u tickvault -f
# =============================================================================

[Unit]
Description=tickvault — Live F&O Trading System
Documentation=https://github.com/SJParthi/tickvault
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=notify
# The app binary. Adjust path to match your deployment.
ExecStart=/opt/tickvault/tickvault
WorkingDirectory=/opt/dlt

# Auto-restart on ANY failure. RestartSec prevents tight restart loops.
Restart=always
RestartSec=3

# Watchdog: app must call sd_notify(WATCHDOG=1) every 30s.
# If no watchdog ping for 60s (2x interval), systemd kills and restarts.
WatchdogSec=60

# Environment
Environment=RUST_LOG=info
Environment=TV_ENVIRONMENT=prod

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tickvault/data /opt/tickvault/config
PrivateTmp=true

# Graceful shutdown: send SIGTERM, wait 30s, then SIGKILL.
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Logging to journald (also captured by Alloy → Loki if configured)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tickvault

[Install]
WantedBy=multi-user.target
