# deploy/systemd/tickvault.service
# =============================================================================
# tickvault — systemd Service Unit
# =============================================================================
# Auto-restarts the trading app on crash. Watchdog kills hung processes.
#
# Install:
#   sudo cp deploy/systemd/tickvault.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tickvault
#   sudo systemctl start tickvault
#
# Check:
#   systemctl status tickvault
#   journalctl -u tickvault -f
# =============================================================================

[Unit]
Description=tickvault — Live F&O Trading System
Documentation=https://github.com/SJParthi/tickvault
After=network-online.target docker.service
Wants=network-online.target docker.service

# Crash-loop safety net (interaction-audit MED — "crash-loop hides under a
# stale gauge"). Restart=always (below) is correct for a genuine one-off crash
# AFTER READY=1, but WITHOUT a start-rate limit a fast crash-loop (e.g. a
# post-READY panic that reproduces on every restart) would restart FOREVER,
# every ~3s, silently. The observability trap: metrics-exporter-prometheus
# re-renders the boot gauges' last value on each scrape while the process is
# briefly alive between restarts, so a crash-loop can keep re-publishing a
# STALE-looking healthy gauge and never surface as a hard failure.
#
# StartLimit* makes that crash-loop OBSERVABLE: > StartLimitBurst restarts
# within StartLimitIntervalSec transitions the unit to `failed` (queryable via
# `systemctl status` / `systemctl is-failed` / an SSM check), which STOPS the
# restart churn and STOPS the stale-gauge re-publish — the app is then genuinely
# dead, its metrics go MISSING, and the market-hours liveness alarm (added in
# this PR) pages.
#
# Values: > 8 restarts in 10 min → `failed`. A healthy slow first boot is ONE
# start attempt that runs long (guarded by TimeoutStartSec=infinity below), NOT
# repeated restarts, so this never trips a legitimate boot. These directives are
# placed in [Unit] (NOT [Service]) per modern systemd — they govern the unit's
# start-rate, unlike Restart=/RestartSec= which live in [Service].
StartLimitIntervalSec=600
StartLimitBurst=8

[Service]
Type=notify
# Run as ec2-user (Amazon Linux 2023 default unprivileged account).
# /opt/tickvault is chowned to ec2-user in the user-data script.
User=ec2-user
Group=ec2-user
# The app binary, dropped by the deploy-aws GitHub Actions workflow.
ExecStart=/opt/tickvault/bin/tickvault
WorkingDirectory=/opt/tickvault

# Self-heal: ensure the QuestDB container EXISTS + is started before the app
# boots. user-data's first-boot `docker compose up -d` runs only once
# (cloud-init), and Docker's `restart: unless-stopped` cannot recreate a
# container that was removed (down/wipe/reset). Without this, a missing QuestDB
# makes the app crash-loop forever on BOOT-02.
#
# This calls scripts/ensure-questdb.sh (idempotent) instead of a bare
# `docker compose ... up -d questdb`, which had THREE bugs (incident 2026-06-08):
#   (a) hosts without the `docker compose` v2 plugin → no-op (docker prints help);
#   (b) the compose SERVICE name is `tv-questdb`, not `questdb`;
#   (c) a recreate needs QDB_PG_* creds from SSM, absent from this unit's env.
# The helper handles running/stopped/removed, tries compose v2→v1→plugin-path,
# and falls back to `docker run` so the box recovers even with no compose.
# Leading `-` = a transient hiccup does NOT hard-fail the unit; the app's own
# 60s `wait_for_questdb_ready` is the real gate.
# NOTE: still does NOT fix QuestDB CRASHING on start (disk-full / volume
# corruption) — see docs/runbooks/aws-docker-daemon-dead.md (`docker logs
# tv-questdb` + `df -h`).
ExecStartPre=-/bin/bash /opt/tickvault/repo/scripts/ensure-questdb.sh

# Boot can legitimately exceed systemd's default 90s TimeoutStartSec: the app
# is Type=notify and only sends READY=1 AFTER its own boot gates complete —
# QuestDB readiness (BOOT-01/BOOT-02), then the infinite-retry daily-universe
# cold build (CSV fetch + NTM resolve + ~767-row lifecycle reconcile). On the
# first-of-day cold boot that whole chain routinely runs past 90s, so the
# default timer would SIGTERM a healthy boot mid-build, producing a crash-loop.
# The app's OWN gates decide boot failure (they HALT + page on a real fault);
# systemd's blunt 90s timer must not pre-empt them. Restart=always (below) still
# restarts a genuine crash that happens AFTER READY=1.
TimeoutStartSec=infinity

# Auto-restart on ANY failure. RestartSec prevents tight restart loops.
# Bounded by StartLimitIntervalSec/StartLimitBurst in [Unit]: after > 8
# restarts in 10 min the unit goes to `failed` and stops restarting (so a true
# crash-loop becomes an observable hard failure instead of an infinite silent
# churn under a re-emitted stale gauge). Recover a `failed` unit manually with
# `systemctl reset-failed tickvault && systemctl start tickvault`.
Restart=always
RestartSec=3

# Watchdog: app must call sd_notify(WATCHDOG=1) every 30s.
# If no watchdog ping for 60s (2x interval), systemd kills and restarts.
WatchdogSec=60

# Environment
Environment=RUST_LOG=info
# TV_ENVIRONMENT selects BOTH config/<env>.toml AND the /tickvault/<env>/ SSM
# prefix (read by boot_helpers::resolve_config_env + secret_manager).
# Single real env (operator 2026-06-30): prod. dev/staging configs were retired;
# base.toml + production.toml are the only env configs. This is STILL the
# no-real-orders data-pull phase — production.toml LOCKS dry_run=true (NO real
# orders) plus a far-future sandbox_only_until belt-and-suspenders. Flipping to
# live trading is a separate, deliberate change (set dry_run=false +
# mode=live AND re-register the static EIP). See config/production.toml.
Environment=TV_ENVIRONMENT=prod
# AWS SDK credentials: pin the default provider chain straight to the EC2
# instance-profile (IMDS) credentials — the ONLY intended credential source on
# the box (SSM reads, and the CloudWatch metrics pipeline that fills the
# Tickvault/Prod namespace).
#
# Root cause (live incident 2026-07-02): this unit hardens with
# ProtectHome=true (below), which masks /home for the service. The AWS SDK's
# default profile-file probe of ~/.aws/config + ~/.aws/credentials
# (HOME=/home/ec2-user under User=ec2-user) therefore fails with EACCES and
# logs "failed to read config file ... Permission denied" on every boot.
# On the pinned Rust SDK (aws-runtime 1.7.5, src/env_config/source.rs — the
# "default config/credentials files MUST NOT return an error" arm) that is a
# warn-and-treat-as-empty, but the probe is pure noise/ambiguity, and every
# CHILD process this unit spawns (e.g. the Groww sidecar's boto3) inherits
# this environment with its own, stricter, profile-file handling. Pointing
# both files at /dev/null removes the unreadable-file probe entirely: the
# chain resolves env -> (empty profile) -> IMDS instance profile
# deterministically, with zero dependency on any on-disk ~/.aws state.
Environment=AWS_CONFIG_FILE=/dev/null
Environment=AWS_SHARED_CREDENTIALS_FILE=/dev/null
# GROWW_MAX_SUBSCRIBE caps the first-run Groww subscription set. Default in code
# is 60 (subset-safe); set to 1000 here to subscribe the FULL ~779-instrument
# Groww universe (NTM-union: indices + F&O underlyings + NIFTY-Total-Market
# constituents). The builder clamps to GROWW_MAX_SUBSCRIPTIONS=1000 (a max, not
# a floor — a prefix-take, indices-first then token-ascending), so 1000 cleanly
# covers the live ~779 count with headroom. Per .claude/rules/project/
# groww-second-feed-scope-2026-06-19.md.
Environment=GROWW_MAX_SUBSCRIBE=1000

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Docker socket exemption (2026-07-13, deploy-hang fix companion):
# ProtectSystem=strict mounts the whole hierarchy read-only (except
# /dev /proc /sys), and connect(2) on a unix socket needs WRITE access to the
# socket inode — so every `docker …` call made by ExecStartPre
# (ensure-questdb.sh) fails with EACCES in the SERVICE-USER context while the
# SAME script works from an SSM/root shell (no sandbox). On a post-nuke boot
# with no QuestDB container/image that meant the self-heal could never bring
# QuestDB up → BOOT-02 halt at 08:30. The `-` prefix ignores a missing socket
# (docker.service not yet up), so this can never fail the unit start; it only
# widens the sandbox by the docker control socket (ec2-user is already in the
# docker group — no privilege change). /var/run is a /run symlink; both listed
# so the exemption holds regardless of which path the CLI resolves.
ReadWritePaths=/opt/tickvault/data /opt/tickvault/config -/var/run/docker.sock -/run/docker.sock
PrivateTmp=true

# Graceful shutdown: send SIGTERM, wait 30s, then SIGKILL.
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Logging to journald (also captured by Alloy → Loki if configured)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tickvault

[Install]
WantedBy=multi-user.target
