# scripts/tv-tunnel/tickvault-tunnel.service
[Unit]
Description=tickvault Tailscale Funnel — exposes the tickvault API (3001) for Claude MCP access
Documentation=file:///home/ec2-user/tickvault/docs/runbooks/claude-mcp-access.md
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple
# Funnel foreground mode so systemd owns the lifecycle.
# Security trim 2026-07-04: ONLY the tickvault API (3001, bearer-auth) is
# funnelled. Prometheus 9090 / Alertmanager 9093 / Grafana 3000 are retired
# (CloudWatch-only migration); QuestDB 9000 is intentionally NO LONGER
# funnelled — it is an auth-less raw-SQL surface (local access only).
ExecStart=/usr/bin/tailscale funnel 3001
ExecStopPost=/usr/bin/tailscale funnel reset
Restart=always
RestartSec=10
# Funnel requires root on Linux to bind low ports on the TS socket.
User=root
Group=root
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
